# root/etc/udev/rules.d/97-nextcloud-kobo.rules
KERNEL=="lo", RUN+="/bin/sh /usr/local/nextcloud-kobo/run.sh"

